# ansible/roles/hazelcast/templates/hazelcast.service.j2
[Unit]
Description=Hazelcast IMDG
After=network.target

[Service]
Type=simple
User=root
ExecStart=/opt/hazelcast/bin/hz start
Restart=on-failure

[Install]
WantedBy=multi-user.target
